# script/wl1251-cal
#!/sbin/openrc-run
description="Extract wl1251 nvs from CAL and write it to /lib/firmware"

depend()
{
        before networking
}

start()
{
        wl1251-extract-nvs
        # probe wl1251 once the nvs is ready
        modprobe wl1251_spi
        eend $?
}
